# speaker-server.service
[Unit]
Description=3D-Speaker Inference Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/speaker
Environment=PYTHONPATH=/root/speaker
ExecStart=/root/miniconda3/envs/3D-Speaker/bin/python server.py
Restart=always
RestartSec=10

# 环境变量
EnvironmentFile=/root/speaker/.env

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=speaker-server

# 安全设置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
